# poker-game.service
[Unit]
Description=Texas Hold'em Poker Game
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/poker-game
Environment=NODE_ENV=production
EnvironmentFile=/opt/poker-game/.env
ExecStart=/usr/bin/node backend/server.js
Restart=always
RestartSec=10
StandardOutput=append:/opt/poker-game/logs/output.log
StandardError=append:/opt/poker-game/logs/error.log

[Install]
WantedBy=multi-user.target
